# vllm-dsv4-head.service
# systemd user unit for the HEAD node (rank 0). Install via install-services.sh.
# Coordination with the worker (and any optional CONFLICTING_SERVICE) is script-level
# (preflight.sh) — systemd deps cannot cross user managers or nodes.
[Unit]
Description=vLLM DeepSeek-V4-Flash-DSpark head (rank 0, TP=2)
StartLimitIntervalSec=3600
StartLimitBurst=6

[Service]
Type=simple
WorkingDirectory=%h/dgx-cluster
# Worst case: 3 boot-dep waits (300s) + peer wait (300s) + worker wait (600s)
TimeoutStartSec=2100
TimeoutStopSec=90
ExecStartPre=/usr/bin/bash %h/dgx-cluster/preflight.sh head
ExecStartPre=/usr/bin/bash %h/dgx-cluster/render-env.sh head
ExecStart=/usr/bin/docker compose --env-file %h/dgx-cluster/.env.dspark -f %h/dgx-cluster/docker-compose.dspark.yml up --exit-code-from vllm-dspark
ExecStop=/usr/bin/docker compose --env-file %h/dgx-cluster/.env.dspark -f %h/dgx-cluster/docker-compose.dspark.yml down -t 30
ExecStopPost=-/usr/bin/docker rm -f vllm-dsv4
Restart=always
RestartSec=20

[Install]
WantedBy=default.target
