# vllm-dsv4-worker.service
# systemd user unit for the WORKER node (rank 1, --headless). Install via install-services.sh.
[Unit]
Description=vLLM DeepSeek-V4-Flash-DSpark worker (rank 1, headless)
StartLimitIntervalSec=3600
StartLimitBurst=10

[Service]
Type=simple
WorkingDirectory=%h/dgx-cluster
TimeoutStartSec=900
TimeoutStopSec=90
ExecStartPre=/usr/bin/bash %h/dgx-cluster/preflight.sh worker
ExecStartPre=/usr/bin/bash %h/dgx-cluster/render-env.sh worker
ExecStart=/usr/bin/docker compose --env-file %h/dgx-cluster/.env.dspark -f %h/dgx-cluster/docker-compose.dspark.yml up --exit-code-from vllm-dspark
ExecStop=/usr/bin/docker compose --env-file %h/dgx-cluster/.env.dspark -f %h/dgx-cluster/docker-compose.dspark.yml down -t 30
ExecStopPost=-/usr/bin/docker rm -f vllm-dsv4
Restart=always
RestartSec=30

[Install]
WantedBy=default.target
