# runtime/vllm-dsv4-watchdog.service
[Unit]
Description=vLLM DSv4 head inference watchdog (restarts head when /health lies)

[Service]
Type=oneshot
ExecStart=/usr/bin/bash %h/dgx-cluster/runtime/watchdog.sh
